# archlinux/icloud.timer
[Unit]
Description=Timer for icloud photo downloader

[Timer]
OnCalendar=*-*-* 1:00:00
Unit=icloud.service

[Install]
WantedBy=timers.target
